# ansible/roles/elgeeko1-scoot-ansible/files/scoot.service
[Unit]
Description=Scoot service
After=network.target

[Service]
User=gimp
Group=gimp
ExecStart=/opt/scoot/run.sh
Restart=always

[Install]
WantedBy=multi-user.target